# templates/upstart.conf
description "Sinopia NPM registry"

start on runlevel [2345]
stop on shutdown
respawn
setuid {{ sinopia_user }}

script
  sinopia -c {{ sinopia_config }}/config.yaml
end script
